# deis/templates/default/container.conf.erb
pre-start script

end script

start on (local-filesystems and net-device-up and runlevel [2345])
stop on runlevel [!2345]

respawn
respawn limit 2 10   # stop if there are 2 respawns in 10 seconds

setuid root
setgid root

script
  exec docker run -p :<%= @port %> -e PORT=<%= @port %> <% if not @slug_dir.nil? %>-v <%= @slug_dir %>/app:/app<% end %> <% @env.each_pair do |k,v| %> -e <%= k.upcase %>='<%= v %>'<% end %> <%= @image %> <% if ! @command.nil? %>/bin/sh -c 'export HOME=/app; cd /app && for profile in `find /app/.profile.d/*.sh -type f`; do . $profile; done && <%= @command %>'<% end %> 2>&1 | logger -p local0.info -t <%= @app_name %>[<%=@c_type%>.<%= @c_num %>]
end script

# hacky workaround for docker port race condition
post-stop script
  while [ "`netstat -an | grep LISTEN | egrep ':<%= @port %>'`" != "" ]; do
    sleep 1
  done
end script
