# dashing/templates/default/upstart-dashboard.conf.erb
#!upstart
# Generated by chef - Do not edit!

description "Dashing Dashboard: <%= @dashboard_name %>"
author      "benbria"

start on runlevel [2345]
stop on runlevel [!2345]
respawn

limit nofile 65536 65536

script
    DASHBOARD_DIR="<%= @dashboard_path %>"
    DASHING_PORT=<%= @dashboard_port %>

    echo $$ > "/var/run/<%= @service_name %>.pid"
    cd "${DASHBOARD_DIR}"
    exec sudo -u <%= node["dashing"]["user"] %> /usr/sbin/dashing-service start -e production -p ${DASHING_PORT}
end script

pre-stop script
    rm "/var/run/<%= @service_name %>.pid"
end script
