# terracotta/templates/default/terracotta-init.erb
#!/bin/sh
#
#
# chkconfig: 2345 55 25
# description: Terracotta L2 Server Daemon                                                                                                                 
#
# Set runlevels to 2 through 5
# Set Startup order to 55 (must be before Tomcat)
# Set Stop order to 25 (must be after Tomcat)
#

# Generated by Chef for <%= node[:fqdn] %>

# where Java is installed
JAVA_HOME="<%= node["java"]["java_home"] %>"

# where tc is installed
TC_HOME="<%= node["terracotta"]["home"] %>"

# where to put stdout/stderr logs
TC_LOGS="$TC_HOME/logs"

# the tc config file
TC_CONFIG="$TC_HOME/tc-config.xml"

# user to run tc as
TC_USER=root
# the identity of this node
TC_SERVER="<%= node[:hostname] %>"
 
start() {
    #
    # Start Terracotta
    #
    echo "Starting Terracotta Server as " ${TC_SERVER}
    su - $TC_USER -c "JAVA_HOME=$JAVA_HOME ${TC_HOME}/bin/start-tc-server.sh -n ${TC_SERVER} -f ${TC_CONFIG} 2>&1 &" > ${TC_LOGS}/terracotta.log &
}
 
stop() {
    #
    # Stop Terracotta
    #
    echo "Stopping Terracotta Server ..."
    su - $TC_USER -c "${TC_HOME}/bin/stop-tc-server.sh 2>&1 &"  > ${TC_LOGS}/terracotta.log &
}
 
case "$1" in
  start)
    start
    ;; 
  stop)
    stop
    ;;
 
  restart)
    # restart Terracotta
    stop
    # wait just a moment to be sure
    sleep 5
    start
    ;;
    
  *)
    echo "Usage $0 start/stop/restart"
    exit 1
esac

